# packaging/apt/postrm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

case "$1" in
    purge)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl disable rginx.service >/dev/null 2>&1 || true
        fi
        rm -f /run/rginx.pid /run/rginx/admin.sock >/dev/null 2>&1 || true
        rm -rf /var/lib/rginx >/dev/null 2>&1 || true
        ;;
esac

exit 0
