# packaging/apt/postinst
#!/bin/sh
set -e

create_rginx_user() {
    if getent passwd rginx >/dev/null 2>&1; then
        return 0
    fi

    shell=/usr/sbin/nologin
    if [ ! -x "$shell" ]; then
        shell=/usr/bin/false
    fi

    adduser \
        --system \
        --group \
        --home /var/lib/rginx \
        --no-create-home \
        --shell "$shell" \
        rginx >/dev/null 2>&1 || true
}

case "$1" in
    configure)
        create_rginx_user
        mkdir -p /var/lib/rginx
        if getent passwd rginx >/dev/null 2>&1; then
            chown rginx:rginx /var/lib/rginx >/dev/null 2>&1 || true
        fi

        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload >/dev/null 2>&1 || true
            systemctl preset rginx.service >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
